# dotnet/Doc/Linux/SupervisorConfigFile/Robot.conf
[program:Robot]
command=dotnet App.Robot.Station.dll
directory=/root/netapp/Robot
environment=ASPNETCORE__ENVIRONMENT=Production
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/root/netapp/supervisor_log/Robot.err.log
stdout_logfile=/root/netapp/supervisor_log/Robot.out.log
